# debian/DEBIAN/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        echo "Configuring pi-clock..."

        # Check if pip3 is available
        if ! command -v pip3 >/dev/null 2>&1; then
            echo "Warning: pip3 not found. Attempting to install python3-pip..."
            apt-get install -y python3-pip 2>/dev/null || {
                echo "Warning: Could not install python3-pip"
                echo "Please install it manually: sudo apt-get install python3-pip"
            }
        fi

        # Install Python packages not available in apt
        if command -v pip3 >/dev/null 2>&1; then
            PIP_PKGS="adafruit-circuitpython-pcd8544 faster-whisper sounddevice yfinance requests"
            echo "Installing Python packages: $PIP_PKGS"
            # Install PortAudio system library required by sounddevice
            apt-get install -y libportaudio2 >/dev/null 2>&1 || true
            pip3 install --break-system-packages --ignore-installed $PIP_PKGS >/dev/null 2>&1 || \
            pip3 install --ignore-installed $PIP_PKGS >/dev/null 2>&1 || {
                echo "Warning: Failed to install some Python packages"
                echo "You may need to install them manually:"
                echo "  pip3 install --break-system-packages $PIP_PKGS"
            }
        else
            echo "Warning: pip3 not available. Please install packages manually:"
            echo "  sudo apt-get install python3-pip"
            echo "  pip3 install --break-system-packages adafruit-circuitpython-pcd8544 vosk sounddevice yfinance requests"
        fi

        # Configure SPI with spi0-0cs overlay so GPIO 8 (CE0) is free for software CS
        echo "Configuring SPI interface..."
        CONFIG_FILE=""
        if [ -f /boot/firmware/config.txt ]; then
            CONFIG_FILE="/boot/firmware/config.txt"
        elif [ -f /boot/config.txt ]; then
            CONFIG_FILE="/boot/config.txt"
        fi

        if [ -n "$CONFIG_FILE" ]; then
            if grep -q "dtparam=spi=on" "$CONFIG_FILE"; then
                sed -i 's/dtparam=spi=on/dtoverlay=spi0-0cs/' "$CONFIG_FILE"
            elif ! grep -q "dtoverlay=spi0-0cs" "$CONFIG_FILE"; then
                echo "dtoverlay=spi0-0cs" >> "$CONFIG_FILE"
            fi
            echo "Note: A reboot is required for SPI changes to take effect."
        else
            echo "Warning: Could not find config.txt. Please add 'dtoverlay=spi0-0cs' manually."
        fi

        # --- Assistant configuration ---
        echo ""
        echo "=== Assistant Configuration ==="
        echo "Which AI backend do you want to use?"
        echo "  1. Ollama (local/network server)"
        echo "  2. Gemini (Google API)"
        echo "  3. Disable Assistant"
        printf "Select [1]: "
        read ASST_CHOICE </dev/tty
        ASST_CHOICE=${ASST_CHOICE:-1}

        SERVICE_FILE="/opt/pi-clock/pi_clock.service"

        if [ "$ASST_CHOICE" = "3" ]; then
            echo "✓ Assistant disabled — menu entry will be hidden"
            sed -i "s|^Environment=ASST_ENABLED=.*|Environment=ASST_ENABLED=0|" "$SERVICE_FILE"
            sed -i "s|^Environment=ASST_BACKEND=.*|Environment=ASST_BACKEND=ollama|" "$SERVICE_FILE"
            sed -i "s|^Environment=OLLAMA_HOST=.*|Environment=OLLAMA_HOST=|" "$SERVICE_FILE"
            sed -i "s|^Environment=OLLAMA_MODEL=.*|Environment=OLLAMA_MODEL=|" "$SERVICE_FILE"
            sed -i "s|^Environment=GEMINI_API_KEY=.*|Environment=GEMINI_API_KEY=|" "$SERVICE_FILE"

        elif [ "$ASST_CHOICE" = "2" ]; then
            printf "Enter your Gemini API key: "
            read GEMINI_KEY </dev/tty
            if [ -z "$GEMINI_KEY" ]; then
                echo "✗ No API key entered. Assistant disabled."
                sed -i "s|^Environment=ASST_ENABLED=.*|Environment=ASST_ENABLED=0|" "$SERVICE_FILE"
            else
                printf "Gemini model [gemini-2.5-flash-lite]: "
                read GEMINI_MDL </dev/tty
                GEMINI_MDL=${GEMINI_MDL:-gemini-2.5-flash-lite}
                echo "✓ Backend  : Gemini"
                echo "✓ Model    : ${GEMINI_MDL}"
                sed -i "s|^Environment=ASST_ENABLED=.*|Environment=ASST_ENABLED=1|" "$SERVICE_FILE"
                sed -i "s|^Environment=ASST_BACKEND=.*|Environment=ASST_BACKEND=gemini|" "$SERVICE_FILE"
                sed -i "s|^Environment=GEMINI_API_KEY=.*|Environment=GEMINI_API_KEY=${GEMINI_KEY}|" "$SERVICE_FILE"
                sed -i "s|^Environment=GEMINI_MODEL=.*|Environment=GEMINI_MODEL=${GEMINI_MDL}|" "$SERVICE_FILE"
                sed -i "s|^Environment=OLLAMA_HOST=.*|Environment=OLLAMA_HOST=|" "$SERVICE_FILE"
                sed -i "s|^Environment=OLLAMA_MODEL=.*|Environment=OLLAMA_MODEL=|" "$SERVICE_FILE"
            fi

        else
            # Ollama (default)
            printf "Enter Ollama server IP or host:port (e.g. 192.168.1.30)\nPress Enter to disable Assistant: "
            read OLLAMA_HOST_INPUT </dev/tty
            if [ -z "$OLLAMA_HOST_INPUT" ]; then
                echo "✓ Assistant disabled — menu entry will be hidden"
                sed -i "s|^Environment=ASST_ENABLED=.*|Environment=ASST_ENABLED=0|" "$SERVICE_FILE"
                sed -i "s|^Environment=OLLAMA_HOST=.*|Environment=OLLAMA_HOST=|" "$SERVICE_FILE"
                sed -i "s|^Environment=OLLAMA_MODEL=.*|Environment=OLLAMA_MODEL=|" "$SERVICE_FILE"
            else
                if echo "$OLLAMA_HOST_INPUT" | grep -q ":"; then
                    OLLAMA_URL="http://${OLLAMA_HOST_INPUT}"
                else
                    OLLAMA_URL="http://${OLLAMA_HOST_INPUT}:11434"
                fi
                echo "Fetching available models from ${OLLAMA_URL}..."
                MODELS_JSON=$(curl -sf --connect-timeout 5 "${OLLAMA_URL}/api/tags" 2>/dev/null || echo "")
                if [ -z "$MODELS_JSON" ]; then
                    echo "Warning: Could not reach Ollama server at ${OLLAMA_URL}."
                    printf "Enter model name manually [deepseek-coder-v2:latest]: "
                    read OLLAMA_MODEL_FINAL </dev/tty
                    OLLAMA_MODEL_FINAL=${OLLAMA_MODEL_FINAL:-deepseek-coder-v2:latest}
                else
                    MODEL_LIST=$(echo "$MODELS_JSON" | python3 -c "
import sys, json
models = [m['name'] for m in json.load(sys.stdin).get('models', [])]
for i, m in enumerate(models, 1):
    print(f'  {i}. {m}')
" 2>/dev/null)
                    MODEL_COUNT=$(echo "$MODELS_JSON" | python3 -c "
import sys, json
print(len(json.load(sys.stdin).get('models', [])))
" 2>/dev/null)
                    if [ -z "$MODEL_LIST" ] || [ "$MODEL_COUNT" = "0" ]; then
                        echo "No models found on server."
                        printf "Enter model name manually [deepseek-coder-v2:latest]: "
                        read OLLAMA_MODEL_FINAL </dev/tty
                        OLLAMA_MODEL_FINAL=${OLLAMA_MODEL_FINAL:-deepseek-coder-v2:latest}
                    else
                        echo ""
                        echo "Available models:"
                        echo "$MODEL_LIST"
                        echo ""
                        printf "Select model number [1]: "
                        read MODEL_SEL </dev/tty
                        MODEL_SEL=${MODEL_SEL:-1}
                        OLLAMA_MODEL_FINAL=$(echo "$MODELS_JSON" | python3 -c "
import sys, json
models = [m['name'] for m in json.load(sys.stdin).get('models', [])]
try:
    print(models[int('${MODEL_SEL}') - 1])
except:
    print(models[0])
" 2>/dev/null)
                    fi
                fi
                echo "✓ Backend  : Ollama"
                echo "✓ Server   : ${OLLAMA_HOST_INPUT}"
                echo "✓ Model    : ${OLLAMA_MODEL_FINAL}"
                sed -i "s|^Environment=ASST_ENABLED=.*|Environment=ASST_ENABLED=1|" "$SERVICE_FILE"
                sed -i "s|^Environment=ASST_BACKEND=.*|Environment=ASST_BACKEND=ollama|" "$SERVICE_FILE"
                sed -i "s|^Environment=OLLAMA_HOST=.*|Environment=OLLAMA_HOST=${OLLAMA_HOST_INPUT}|" "$SERVICE_FILE"
                sed -i "s|^Environment=OLLAMA_MODEL=.*|Environment=OLLAMA_MODEL=${OLLAMA_MODEL_FINAL}|" "$SERVICE_FILE"
                sed -i "s|^Environment=GEMINI_API_KEY=.*|Environment=GEMINI_API_KEY=|" "$SERVICE_FILE"
            fi
        fi

        # Create symbolic link for systemd service
        if [ ! -f /etc/systemd/system/pi_clock.service ]; then
            ln -s /opt/pi-clock/pi_clock.service /etc/systemd/system/pi_clock.service
        fi

        # Reload systemd daemon
        systemctl daemon-reload

        # Enable and start the service
        echo "Enabling pi-clock service..."
        systemctl enable pi_clock.service >/dev/null 2>&1 || true

        echo "Starting pi-clock service..."
        systemctl start pi_clock.service >/dev/null 2>&1 || {
            echo "Warning: Service failed to start. Check logs with: journalctl -u pi_clock.service -n 50"
        }

        # Check if service is running
        if systemctl is-active --quiet pi_clock.service; then
            echo "✓ pi-clock installed and running successfully!"
            echo ""
            echo "Service management commands:"
            echo "  Status:  sudo systemctl status pi_clock.service"
            echo "  Stop:    sudo systemctl stop pi_clock.service"
            echo "  Start:   sudo systemctl start pi_clock.service"
            echo "  Restart: sudo systemctl restart pi_clock.service"
            echo "  Logs:    journalctl -u pi_clock.service -f"
        else
            echo "✗ Service not running. Check hardware connections and logs."
        fi
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
